# deploy/systemd/tone-of-voice-telegram-bot.service.example
[Unit]
Description=tone-of-voice Telegram drafting bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/tone-of-voice
EnvironmentFile=/opt/tone-of-voice/.env
ExecStart=/opt/tone-of-voice/.venv/bin/python scripts/run_telegram_bot.py --output-dir /opt/tone-of-voice/data/bot
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
